# root/etc/services.d/huedream/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-022}

umask "$UMASK_SET"

cd /opt/huedream || exit

exec \
	HueDream
